# profile.d/VHI3-settings.sh
setxkbmap -option caps:swapescape
setxkbmap -layout us,de,ir
setxkbmap -option 'grp:alt_shift_toggle'
export PATH="~/.local/bin:$PATH"
#export PATH="/usr/local/texlive/2022/bin/x86_64-linux:$PATH"
#export INFOPATH=$INFOPATH:/usr/local/texlive/2022/texmf-dist/doc/info
#export MANPATH=$MANPATH:/usr/local/texlive/2022/texmf-dist/doc/man
